# asterisk.conf
[program:festival]
command=/usr/bin/festival --server

[program:asterisk]
command=/usr/sbin/asterisk -T -W -U asterisk -p -nvvvdddf